# dev-utils/resources/supervisor/conf.d/sshd.conf
[program:sshd]
directory=/usr/sbin/
command=/usr/sbin/sshd -D -e
