# docker/dbdozer/supervisor/dbdozer.conf
[program: dbdozer]
directory = /opt/dbdozer
command = /opt/dbdozer/start-dbdozer.sh
autorestart = true
