# docker/dbdozer/supervisor/lb-blacklist-generator.conf
[program: lb-blacklist_generator]
directory = /opt/dbdozer
command = /opt/dbdozer/start-lb-blacklist-generator.sh
