# integrations/supervisor/apps-conf/integrations-file.conf
[program: integrations-file]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh file /opt/bios/configuration/integrations/file.ini
