# integrations/supervisor/apps-conf/integrations-s3.conf
[program: integrations-s3]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh s3 /opt/bios/configuration/integrations/s3.ini
