# integrations/supervisor/apps-conf/integrations-mongodb.conf
[program: integrations-mongodb]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh mongodb /opt/bios/configuration/integrations/mongodb.properties
