# integrations/supervisor/apps-conf/integrations-mysql.conf
[program: integrations-mysql]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh mysql /opt/bios/configuration/integrations/mysql.properties
