# integrations/supervisor/apps-conf/integrations-kafka.conf
[program: integrations-kafka]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh kafka /opt/bios/configuration/integrations/kafka.ini
