# integrations/supervisor/apps-conf/integrations-webhook.conf
[program: integrations-webhook]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh webhook /opt/bios/configuration/integrations/webhook.ini
autorestart = true
