# integrations/supervisor/apps-conf/integrations-rest-client.conf
[program: integrations-rest-client]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh rest-client /opt/bios/configuration/integrations/rest_client.ini
