# integrations/supervisor/apps-conf/integrations-sql-pull.conf
[program: integrations-sql-pull]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh sql-pull /opt/bios/configuration/integrations/sql_pull.ini
