# integrations/supervisor/apps-conf/integrations-file-tailer.conf
[program: integrations-file-tailer]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh file-tailer /opt/bios/configuration/integrations/file-tailer.ini
