# integrations/supervisor/apps-conf/integrations-postgres.conf
[program: integrations-postgres]
directory = /opt/bios
command = /opt/bios/apps/integrations/start-bios-integrations.sh postgres /opt/bios/configuration/integrations/postgres.properties
