# sql/resources/supervisor/conf.d/trino.conf
[program: trino]
directory = /opt/bios/server-apps/trino
command = /opt/bios/server-apps/trino/bin/launcher --etc-dir=/opt/bios/configuration/trino run
